# examples/landshark_demo/radl/files/kestrel/ocu/etc/network/if-pre-up.d/iptablesload
#!/bin/sh
iptables-restore < /opt/hacms/etc/iptables.rules
exit 0
